# deploy/gaottt-embedder.service
[Unit]
Description=GaOTTT Embedding Service (RURI shared host model)
After=network.target

[Service]
Type=simple
User=gaottt
WorkingDirectory=/opt/gaottt
Environment=PATH=/opt/gaottt/.venv/bin
ExecStart=/opt/gaottt/.venv/bin/python -m gaottt.embedding.service --host 127.0.0.1 --port 7879
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
